# adds/etc/supervisor.d/docker-service.conf
[program:dockerd]
command=dockerd --host=unix:///var/run/docker.sock --host=tcp://0.0.0.0:2375
stderr_logfile = /var/log/supervisord/dockerd.err
stdout_logfile = /var/log/supervisord/dockerd.log
